# nginx/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
stderr_logfile=/var/log/supervisor/nginx-err.log
stdout_logfile=/var/log/supervisor/nginx-out.log
priority=10

[program:cert-watcher]
command=/opt/scripts/watch-certs.sh
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
stderr_logfile=/var/log/supervisor/cert-watcher-err.log
stdout_logfile=/var/log/supervisor/cert-watcher-out.log
priority=20 